# packaging/systemd/vaydeer-studio.service
[Unit]
Description=Vaydeer Studio read-only JP-1011 HID keepalive and binding service
After=graphical-session.target
Wants=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/vaydeer-studiod --log-level info
Restart=on-failure
RestartSec=3
Environment=PYTHONUNBUFFERED=1
Environment=VAYDEER_STUDIO_LOG_LEVEL=INFO
NoNewPrivileges=true
PrivateTmp=true
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
